# Build_Web_With_Flask/Building web applications with Flask_Code/chapter10/chapter10/project.conf
description "uWSGI application my project"

start on runlevel [2345]
stop on runlevel [!2345]

setuid your-user
setgid www-data

exec /path/to/uwsgi --ini /path/to/ini/file.ini